# templates/default/etc/init.d/firewall.erb
#!/bin/sh

###
# Generated by Chef for [<%= node[:fqdn] %>]
# Environment: <%= node.chef_environment %>
###

# Firewall start/stop script.

# Comments to support LSB init script conventions
### BEGIN INIT INFO
# Provides: firewall
# Required-Start: $local_fs $network $remote_fs
# Should-Start: ntpd
# Required-Stop: $local_fs $network $remote_fs
# Default-Start:  2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: start and stop firewall
# Description: Firewall service.
### END INIT INFO

PATH=/sbin:/bin:/usr/sbin:/usr/bin
. /lib/lsb/init-functions

NAME=firewall
DESC="Firewall service"
RULES=/etc/iptables.rules

case "$1" in
    start)
    log_daemon_msg "Starting $DESC" "$NAME"
    if test -f $RULES
    then
        /bin/bash $RULES start > /dev/null 2>&1
        RET=$?
        log_end_msg $RET
        exit $RET
    else
        log_failure_msg "Couldn't find firewall configuration ($RULES)"
        exit 1
    fi
    ;;
    stop)
    log_daemon_msg "Stopping $DESC" "$NAME"
    if test -f $RULES
    then
        /bin/bash $RULES stop > /dev/null 2>&1
        log_end_msg $?
    else
        log_failure_msg "Couldn't find firewall configuration ($RULES)"
        exit 1
    fi
    ;;
    restart|reload)
        $0 stop
        $0 start
    ;;
    *)
      basename=`basename "$0"`
      echo "Usage: $basename  {start|stop|restart|reload}"
      exit 1
    ;;
esac

exit 0
